# release/sutra-engine/examples/sutra-engine.service
[Unit]
Description=Sutra Engine - High Performance Memory Storage
After=network.target

[Service]
Type=simple
User=sutra
Group=sutra
# Adjust paths as necessary
ExecStart=/usr/local/bin/sutra-engine
WorkingDirectory=/var/lib/sutra-engine
Restart=always
RestartSec=3

# Environment Configuration
Environment="STORAGE_PORT=50051"
Environment="STORAGE_PATH=/var/lib/sutra-engine/data"
Environment="RUST_LOG=info"
# Environment="SUTRA_SECURE_MODE=true"

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
