# release/sutra-engine/examples/sutra-engine.service
[Unit]
Description=Sutra Engine Standalone
After=network.target

[Service]
Type=simple
User=sutra
WorkingDirectory=/opt/sutra-engine
ExecStart=/opt/sutra-engine/sutra-engine
Restart=always
RestartSec=10

# Environment variables
Environment=STORAGE_PORT=50051
Environment=STORAGE_PATH=/var/lib/sutra-engine
# Environment=SUTRA_SECURE_MODE=true
# Environment=SUTRA_AUTH_SECRET=your-32-char-secret

[Install]
WantedBy=multi-user.target
